# build/iso-profile/airootfs/etc/systemd/system/hatan-wifi-autoconnect.service
[Unit]
Description=HATAN OS default WiFi autoconnect
After=NetworkManager.service iwd.service
Wants=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hatan-wifi-autoconnect.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
